# autostart/fcitx5.desktop
[Desktop Entry]
Type=Application
Exec=fcitx5
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[en_US]=fcitx5
Name=fcitx5
Comment[en_US]=Start Fcitx5 input method framework
Comment=Start Fcitx5 input method framework
